# root/etc/services.d/gollum/run
#!/command/with-contenv sh

export HOME="/config"
/command/s6-applyuidgid -g ${PGID:-911} -u ${PUID:-911} /config/run.sh
